# etc/supervisord.d/backup_site.ini
[program:backup_site]
command=/root/bin/backup_site.sh
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

