# etc/supervisord.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
priority=10
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

