# etc/supervisord.d/setup_collabtive.ini
[program:setup_collabtive]
command=/root/bin/setup_collabtive.sh
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

